# ansible/roles/common/templates/disk-cleanup.service.j2
[Unit]
Description=Disk cleanup

[Service]
Type=oneshot
ExecStart=/usr/local/bin/disk-cleanup.sh
